# com.financial.dashboard.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.financial.dashboard.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/matiasmollerv/Documents/Claude/FinancialDashboard/run_weekly_update.sh</string>
    </array>

    <!-- Ejecutar cada lunes a las 07:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/matiasmollerv/Documents/Claude/FinancialDashboard/logs/launchd_out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/matiasmollerv/Documents/Claude/FinancialDashboard/logs/launchd_err.log</string>

    <!-- Reintentar si el Mac estaba durmiendo -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
